# overlay/common/etc/systemd/system/openrig-display-boot.service
[Unit]
Description=openRig — Boot display message
DefaultDependencies=no
After=local-fs.target
Before=sysinit.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/openrig-display "Booting..."
RemainAfterExit=no

[Install]
WantedBy=sysinit.target
