# overlay/common/etc/systemd/system/openrig-hosts-update.timer
[Unit]
Description=openRig — Daily YSF/FCS host list refresh

[Timer]
OnBootSec=2min
OnUnitActiveSec=24h

[Install]
WantedBy=timers.target
